# packaging/rpi/systemd/openarcade-configd.service
[Unit]
Description=OpenArcade USB serial configuration service
Requires=openarcade-gadget-manager.service
After=openarcade-gadget-manager.service

[Service]
Type=simple
EnvironmentFile=-/etc/openarcade/openarcade.env
WorkingDirectory=/opt/openarcade/app/server
ExecStart=/opt/openarcade/venv/bin/python /opt/openarcade/app/server/serial_config_service.py --device /dev/ttyGS0 --config $OPENARCADE_CONFIG_PATH
Restart=on-failure
RestartSec=3
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
